# systemd/logger.service
[Unit]
Description=Logger Script Service

[Service]
ExecStart=/usr/local/bin/logger
